# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/sshd.err.log
stdout_logfile=/var/log/supervisor/sshd.out.log
priority=999

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
priority=998

[program:spring-boot]
command=java -Xmx512m -Xms256m -jar /app/demo.jar
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/spring-boot.err.log
stdout_logfile=/var/log/supervisor/spring-boot.out.log
priority=997
stopasgroup=true
killasgroup=true
startsecs=30
startretries=3
